# root_alpine_pg16/etc/services.d/postgres/run
#!/usr/bin/with-contenv sh

echo "Starting postgres..."
s6-setuidgid postgres /usr/bin/postgres16 -D /config/postgres/
